# systemd_service_files/ubrp_persistence.service
[Unit]
Description=Persist UBPR rules across restarts.

[Service]
Type=simple
ExecStart=/usr/bin/UBRP -q

[Install]
WantedBy=multi-user.target

